# srcs/nginx/supervisord.conf
[supervisord]
user     = root
nodaemon = true

[program:sshd]
command  = /usr/sbin/sshd

[program:nginx]
command  = nginx -g 'daemon off;'
