# srcs/phpmyadmin/supervisord.conf
[supervisord]
user        = root
nodaemon    = true

[program:php-fpm7]
command     = php-fpm7
autoRestart = true

[program:nginx]
command     = nginx -g 'daemon off;'
autoRestart = true